# infra/systemd/mediamtx.service
[Unit]
Description=MediaMTX — RTSP/HLS/WebRTC media server for Heldar
Documentation=https://github.com/bluenviron/mediamtx
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=heldar
Group=heldar
ExecStart=/usr/local/bin/mediamtx /etc/heldar/mediamtx.yml
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
